# baseunits/job-inbox@.path
# Path to container inbox

[Unit]
Description=Watch for a container inbox for job files to appear there

[Path]
DirectoryNotEmpty=/var/vmalwc/boxes/%f/inbox
